# appvm-scripts/etc/xdgautostart/qubes-qrexec-fork-server.desktop
[Desktop Entry]
Name=qrexec-fork-server
Icon=qubes
Exec=/usr/bin/qrexec-fork-server
Terminal=false
Type=Application
