# appvm-scripts/lib/systemd/user/pipewire.service.d/30_qubes.conf
[Unit]
# Avoid pipewire in root's session to not steal user's connection to the audiovm
ConditionUser=!root
