# appvm-scripts/lib/systemd/user/pipewire.service.d/30_qubes.conf
[Unit]
# Only enable if PipeWire service is turned on,
# or if PulseAudio is not installed.  This ensures
# that PipeWire and PulseAudio do not conflict.
ConditionPathExists=|/run/qubes-service/pipewire
ConditionPathExists=|!/usr/bin/pulseaudio
