# appvm-scripts/lib/systemd/user/qubes-graphical-session.target
[Unit]
Description=Qubes graphical session
BindsTo=graphical-session.target
