# appvm-scripts/etc/profile.d/qubes-gui.csh
setenv DISPLAY ":0"
setenv _JAVA_AWT_WM_NONREPARENTING "1"
