# appvm-scripts/etc/profile.d/qubes-session.sh
if [ -z "$QUBES_ENV_SOURCED" ]; then
    . /tmp/qubes-session-env
fi
