# appvm-scripts/etc/profile.d/qubes-gui.sh
export DISPLAY=:0 _JAVA_AWT_WM_NONREPARENTING=1
