# appvm-scripts/etc/systemd/system/lightdm.service.d/qubes-guivm-vnc.conf
[Unit]
ConditionPathExists=/run/qubes-service/lightdm

[Service]
ExecStartPre=/usr/bin/qubes-run-x11vnc pre
ExecStartPost=/usr/bin/qubes-run-x11vnc
ExecStopPost=/usr/bin/qubes-run-x11vnc post
