# pulse/pulseaudio.service.d/30_qubes.conf
[Unit]
# Only enable if PipeWire service is turned off,
# or if PipeWire is not installed.  This ensures
# that PipeWire and PulseAudio do not conflict.
ConditionPathExists=|!/run/qubes-service/pipewire
ConditionPathExists=|!/usr/bin/pipewire

# Only enable if the Qubes configuration file is available
ConditionPathExists=/etc/pulse/qubes-default.pa

[Service]
ExecStartPre=-/usr/bin/qubesdb-read -w /qubes-audio-domain-xid
ExecStart=
ExecStart=/usr/bin/pulseaudio --start -n --file=/etc/pulse/qubes-default.pa --exit-idle-time=-1 --daemonize=no --log-target=journal
